# providers/swarm/scripts/lib/kutab-autoscaler.timer
[Unit]
Description=Run the Kutab autoscaler every 2 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=2min
AccuracySec=15s

[Install]
WantedBy=timers.target
